# systemd/sanity-node-public-rehearsal.service
[Unit]
Description=Generate isolated Sanity Node public-mode rehearsal dashboard
Documentation=https://github.com/vonH23rz/sanity-node
After=network-online.target
Wants=network-online.target
ConditionPathExists=/opt/sanity-node-public-rehearsal/config/config.yaml
ConditionPathExists=/opt/sanity-node-public-rehearsal/app/scripts/run-public-rehearsal.sh

[Service]
Type=oneshot
User=sanity-node
Group=sanity-node
WorkingDirectory=/opt/sanity-node-public-rehearsal/app

Environment=SANITY_NODE_REHEARSAL_ROOT=/opt/sanity-node-public-rehearsal
Environment=SANITY_NODE_APP_ROOT=/opt/sanity-node-public-rehearsal/app
Environment=SANITY_NODE_CONFIG=/opt/sanity-node-public-rehearsal/config/config.yaml
Environment=SANITY_NODE_OUTPUT=/opt/sanity-node-public-rehearsal/html/index.html
Environment=SANITY_NODE_LOG=/opt/sanity-node-public-rehearsal/logs/generator.log
Environment=SANITY_NODE_PROTECTED_ROOT=/opt/homelab-dashboard
EnvironmentFile=-/etc/default/sanity-node-public-rehearsal

ExecStart=/opt/sanity-node-public-rehearsal/app/scripts/run-public-rehearsal.sh
UMask=0002

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadOnlyPaths=/opt/homelab-dashboard
ReadWritePaths=/opt/sanity-node-public-rehearsal
RestrictSUIDSGID=true
LockPersonality=true
